# dist/DEBIAN/postinst
#!/bin/sh
systemctl enable fxdmx.service
systemctl start fxdmx.service