# dist/DEBIAN/prerm
#!/bin/sh
systemctl stop fxdmx.service
systemctl disable fxdmx.service